# .config/systemd/user/sommelier-x@0.service.d/override.conf
[Service]
Environment="SOMMELIER_FRAME_COLOR=#2E3440"
Environment="SOMMELIER_ACCELERATORS=Super_L,<Alt>bracketright,<Alt>bracketleft,<Alt>minus,<Alt>equal"
